# monthly-revenue-scheduler.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.twse.monthly-revenue-fetch</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/liangfuting/Desktop/台股數據資料抓取＿桌機/.venv/bin/python</string>
        <string>/Users/liangfuting/Desktop/台股數據資料抓取＿桌機/scripts/fetch_monthly_revenue_scheduled.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>/Users/liangfuting/Desktop/台股數據資料抓取＿桌機</string>

    <key>StandardOutPath</key>
    <string>/Users/liangfuting/Desktop/台股數據資料抓取＿桌機/monthly_revenue_launchd_stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/liangfuting/Desktop/台股數據資料抓取＿桌機/monthly_revenue_launchd_stderr.log</string>

    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Day</key>
            <integer>10</integer>
            <key>Hour</key>
            <integer>9</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Day</key>
            <integer>11</integer>
            <key>Hour</key>
            <integer>9</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
